# webhook-server.service
[Unit]
Description=GitHub Webhook Server for AIXBiz
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/aixbiz
Environment="NODE_ENV=production"
Environment="GITHUB_WEBHOOK_SECRET=your-secret-here"
Environment="WEBHOOK_PORT=9000"
ExecStart=/usr/bin/node /var/www/aixbiz/webhook-server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target